# config/supervisor/frpc.ini

[program:frpc]

command=/opt/hit-star-gateway/frpc -c /opt/hit-star-gateway/frpc.toml
autostart=true
autorestart=true

stderr_logfile=/var/log/frpc.err.log
stdout_logfile=/var/log/frpc.out.log