# config/supervisor/mihomo.ini
[program:mihomo]
command=/mihomo
autostart=true
autorestart=true

stderr_logfile=/var/log/mihomo.err.log
stdout_logfile=/var/log/mihomo.out.log
