# config/supervisor/sshd.ini
[program:sshd]
command=/usr/sbin/sshd -D -e
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/sshd.err.log
stdout_logfile=/var/log/supervisor/sshd.out.log
